# supervisord.ini
[supervisord]
nodaemon=true

[program:puma]
command=bundle exec puma  -p 3000 -C config/puma.rb
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sidekiq_worker]
command= bundle exec sidekiq -C config/sidekiq.yml
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

stopsignal=QUIT
